# dev/filesystem/etc/profile.d/nix-nice-dev.sh
#
# additional config/setup to perform at login
#
# Note: Linux's /etc/profile will source all scripts ending in ".sh"
# under /etc/profile.d, so put this script there
#
#


# Configure the git bash prompt
GIT_PROMPT_ONLY_IN_REPO=1;
GIT_PROMPT_THEME=Solarized;
